# Service_Files/ram-runner.service
[Unit]
Description=RAM Runner (USB preload → /dev/shm with Sigilbook isolation)
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ram_runner
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=yes
MemoryDenyWriteExecute=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
